# config/shadowbroker.conf
[program:shadowbroker]
command=python shadowbroker.py
numprocs=1
directory=/home/pi/proj/shadowbroker
umask=022
priority=999
autostart=true
autorestart=true